# supervisord.conf
[program:celeryd]
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py celeryd -E -B --loglevel=DEBUG

[program:celerycam]
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py celerycam
